# srs/files/etc/cron.daily/cleanup-srs-thumbs
#!/bin/sh
find /var/www/static/live/thumbs -type f -name "*.jpg" -mtime +7 -print -delete
