# srs/files/etc/systemd/system/rc-hls-flush.service
[Unit]
Description=RC HLS directory flush (ts/m3u8 cleanup)
After=network.target

[Service]
Type=oneshot
Environment=HLS_DIR=/var/www/static/live/hls
Environment=HLSFIX_DIR=/var/www/static/live/hlsfix/live
# legacy/flat output path (some old writers drop .ts directly under hlsfix/)
Environment=HLSFIX_ROOT=/var/www/static/live/hlsfix
Environment=LOG_FILE=/var/log/rc-hls-flush.log
# 기본: 3시간 이상 된 파일만 삭제 (SRS 자체 hls_cleanup + 백업 정리용)
Environment=RETENTION_MINUTES=180
# tmp/part 임시파일은 10분 이상이면 필요없음
Environment=TMP_RETENTION_MINUTES=10
# 긴급(디스크 부족): 10분 이상 된 파일 삭제 (SRS 관련 최우선 정리)
Environment=EMERGENCY_RETENTION_MINUTES=10
Environment=DISK_PCT_EMERGENCY=90
# 최소 여유공간 1GB 확보를 최우선으로 트리거
Environment=DISK_MIN_FREE_MB=1024
Environment=LARGE_TS_MB=100
ExecStart=/usr/local/bin/rc-hls-flush.sh
Nice=10

