# srs/files/etc/systemd/system/rc-rtmp-monitor.service
[Unit]
Description=RC RTMP Monitor (periodic snapshot)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rc-rtmp-monitor.sh
